# deploy/macos/com.ccbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ccbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/ABSOLUTE/PATH/TO/ccbot</string>
        <string>ccbot</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/ccbot</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/.ccbot/ccbot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/.ccbot/ccbot.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Include ~/.local/bin so the Claude Code native installer binary
             (~/.local/bin/claude) is resolvable by subprocess calls from
             ccbot (e.g. `claude --version` in update_watcher). -->
        <key>PATH</key>
        <string>/Users/USERNAME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
